# oc-patches/cc1-app/misc-app-patch/rc.local
if [ -e /app/noapp.kill ]; then
  echo "Kill file noapp.kill present! Continuing boot without /app/app..."
  rm -f /app/noapp.kill 2>/dev/null
  echo "/app/noapp.kill deleted."
  exit 0
fi
