# oc-patches/setup/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

if [ -f /app/app ]; then
    /app/force-upgrade
fi
