# oc-patches/services/install-extra-dependencies-service/oc-install-extra-packages
#!/bin/sh /etc/rc.common

START=82

boot() {
    cd /app/packages
    opkg install --cache . *.ipk
    pip install pyserial-3.4-py2.py3-none-any.whl
}
